# atsc.service
[Unit]
Description=Actuated Traffic Signal Controller (ATSC)
StartLimitBurst=5
StartLimitIntervalSec=80s

[Service]
Type=forking

TimeoutSec=10
Restart=on-abnormal
RestartSec=16

WorkingDirectory=/home/pi/atsc/
ExecStart=/home/pi/atsc/run_pi.sh
PIDFile=/home/pi/atsc/atsc.pid

KillMode=mixed
Restart=on-abnormal
RestartSec=15

SuccessExitStatus=1
RestartPreventExitStatus=1

[Install]
WantedBy=default.target
