# atsc.service
[Unit]
Description=Actuated Traffic Signal Controller (ATSC)

[Service]
Type=forking

WorkingDirectory=%h/atsc-async/
ExecStart=/bin/sh %h/atsc-async/start.sh
ExecStop=/bin/sh %h/atsc-async/stop.sh

Restart=on-success
RestartSec=1
TimeoutStartSec=15
TimeoutStopSec=5

MemoryMax=64M
CPUQuota=50%

[Install]
WantedBy=default.target
